# templates/default/hiveserver2.service.erb
[Unit]
Description = Hive metaserver
After = syslog.target network.target remote-fs.target <%= @deps %>
StartLimitIntervalSec=120
StartLimitBurst=10

[Service]
User = <%= node['hive2']['user'] %>

ExecStartPre=/bin/bash -c 'while ! /usr/bin/host <%= @hive_metastore_endpoint %> ; do if systemctl list-units --full -all | grep -Fq 'hivemetastore'; then if systemctl is-active --quiet hivemetastore; then /bin/sleep 5; else ! break; fi; else /bin/sleep 5; fi;  done'
ExecStart = <%= node['hive2']['base_dir'] %>/bin/hiveserver2
LimitNOFILE = 20000
TimeoutSec= 10min
Restart=on-failure
RestartSec=2

Type = simple

[Install]
WantedBy = multi-user.target
